# roles/vrrp/templates/origin-master-api.service.d/origin-master-api.service
[Unit]
Description=Atomic OpenShift Master API

[Service]
TimeoutStartSec=300
